# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node server.js

[program:python]
command=python server.py